# files/system/etc/systemd/user/swayzzite-gtk-setup.service
[Unit]
Description=Swayzzite First-Login GTK Theme Setup
ConditionPathExists=!%h/.config/swayzzite-gtk-setup-done
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/swayzzite-gtk-setup
RemainAfterExit=yes

[Install]
WantedBy=default.target
